# etc/supervisor.d/agent.ini

[program:php-fpm]
priority=10
directory=/
command=/usr/bin/php-fpm -R -F
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/phpfpm.log
redirect_stderr=true

[program:nginx]
priority=15
directory=/
command=/usr/sbin/nginx
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
redirect_stderr=true

